# files/srv-sites.mount
# /etc/systemd/system/srv-sites.mount
# Managed by Ansible - DO NOT EDIT BY HAND

[Unit]
Description=/srv/sites (read-only, mount)
Requires=zfs.target
After=zfs.target

[Mount]
What=none
Where=/srv/sites
Type=aufs
Options=rw,br=/etc/skel/sites:/home/rob/sites:/home/julie/sites

[Install]
WantedBy=multi-user.target
